# build/linux/wlagent/wlagent.service
[Unit]
Description=wlagent

[Service]
Type=simple

WorkingDirectory=/opt/workload-agent/bin/
ExecStart=/opt/workload-agent/bin/wlagent rungrpcservice
TimeoutStartSec=0
Restart=on-failure
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target
